# cf-workers-sub.service
[Unit]
Description=CF-Workers-SUB Node.js Service
Documentation=https://github.com/cmliu/CF-Workers-SUB
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/cf-workers-sub
Environment=NODE_ENV=production
ExecStart=/usr/bin/node app.js
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
Restart=always
RestartSec=10

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cf-workers-sub

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cf-workers-sub/data /opt/cf-workers-sub/logs

[Install]
WantedBy=multi-user.target
